# src/init/hermesd-openrc
#!/sbin/openrc-run

name="hermesd"
description="System Upgrade Watchdog Service"

command="/usr/libexec/hermesd.py"

pidfile="/run/${RC_SVCNAME}.pid"
logfile="/var/log/hermes-daemon.log"

depend() {
    after net
    need dbus
    use logger
}

start_pre() {
    checkpath --directory --mode 0755 /run
    checkpath --directory --mode 0755 /var/log
    touch "${logfile}"
}

start() {
    ebegin "Starting ${name}"
    export HERMESD_LOGFILE="${logfile}"
    start-stop-daemon --start --background --make-pidfile --pidfile "${pidfile}" \
        --exec "${command}" -- ${command_args} >> "${logfile}" 2>&1
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    start-stop-daemon --stop --pidfile "${pidfile}" --retry 5
    eend $?
}
